# app/docker/conf.d/appium.conf
[program:appium]
command=appium server -p %(ENV_APPIUM_PORT)s --log-timestamp --log-no-colors --allow-cors --allow-insecure=adb_shell
autorestart=true
